# deploy/templates/mem0-mcp.conf.j2
[program:mem0-mcp]
command={{ mem0_home }}/venv/bin/python {{ mem0_home }}/server/run_mcp_server.py
directory={{ mem0_home }}
user={{ mem0_user }}
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile={{ mem0_home }}/logs/mcp.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PYTHONPATH="{{ mem0_home }}",MEM0_CONFIG="{{ mem0_home }}/config/mem0.config.json"
