# scripts/automation/com.ninaivalaigal.graph-test.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ninaivalaigal.graph-test</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/Caskroom/miniconda/base/envs/nina/bin/python</string>
        <string>/Users/swami/Downloads/simple_enhanced_graph_test.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/swami/WorkSpace/ninaivalaigal</string>

    <key>StartInterval</key>
    <integer>600</integer> <!-- every 10 min -->

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/graph-test.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/graph-test-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/Caskroom/miniconda/base/envs/nina/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
